# files-nexus-pro-nx-6/etc/rc.local
. /etc/backone-release

# Netify NFT Rule to block
nft -f /etc/nftables.d/10-netify.nft

echo $BACKONE_MODEL > /tmp/sysinfo/model
echo $BACKONE_BOARD_NAME > /tmp/sysinfo/board_name

/usr/bin/config-rtty.sh

exit 0
